# supervisord.conf
[supervisord]
logfile=/dev/null
pidfile=/var/run/supervisord.pid
nodaemon=true

[program:bgl_hubot]
command=/hubot/bin/hubot --adapter slack
directory=/hubot
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
auto_start=true
autorestart=true
